# crontab/nutch-tab
SHELL=/bin/bash 
0 */6 * * * root . /root/env.sh && $NUTCH_HOME/crawler.sh www-openstack 1 $NUTCH_HOME/urls/www $DEFAULT_DEPTH https://www.openstack.org $DEFAULT_TOP >> /var/log/nutch_cron_www.log 2>&1
0 */6 * * * root . /root/env.sh && $NUTCH_HOME/crawler.sh docs-openstack 2 $NUTCH_HOME/urls/docs $DEFAULT_DEPTH https:/docs.openstack $DEFAULT_TOP >> /var/log/nutch_cron_docs.log 2>&1
0 */6 * * * root . /root/env.sh && $NUTCH_HOME/crawler.sh superuser-openstack 3 $NUTCH_HOME/urls/superuser $DEFAULT_DEPTH "" $DEFAULT_TOP >> /var/log/nutch_cron_superuser.log 2>&1
