# crontab/nutch-tab
SHELL=/bin/bash 
*/5 * * * * root . /root/env.sh && $NUTCH_HOME/crawler.sh www-openstack 1 $NUTCH_HOME/urls/www 10 https://www.openstack.org 100 >> /var/log/nutch_cron.log 2>&1
